# packaging/sprout.service
[Unit]
Description=Sprout hardware-info probe for LettuceAI
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=@BIN@
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
